# conf/supervisor.conf
[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf

[program:mongodb]
command=/usr/bin/mongod --config /etc/mongodb.conf

[program:sshd]
command=/usr/sbin/sshd -D
